# crontabs/athena.cron
MAILTO='william+cron@tuffbizz.com'

0 0 * * * /home/william/scripts/cron-shunt /home/william/scripts/wwwbackup
0,30 * * * * /home/william/scripts/cron-shunt /home/william/scripts/twitter-fortune-bot.pl
