# crontabs/pan.cron
MAILTO='william+cron@tuffbizz.com'

0 0 * * * /home/william/scripts/cron-shunt /home/william/scripts/music-tree
0 0 * * * /home/william/scripts/cron-shunt /home/william/scripts/music-clean
0 12 * * * /home/william/scripts/cron-shunt /home/william/scripts/dailymail.rb
0 */12 * * * /home/william/scripts/cron-shunt /home/william/experiments/points/points.sh
