# crontabs/pan.cron
MAILTO='william+cron@tuffbizz.com'

0 0 * * * /home/william/scripts/cron-shunt /home/william/scripts/music-tree
0 0 * * * /home/william/scripts/cron-shunt /home/william/scripts/music-clean
0 1 * * * /home/william/scripts/cron-shunt /home/william/scripts/music-backup
0 12 * * * /home/william/scripts/cron-shunt /home/william/scripts/dailymail.rb
