# scripts/crontabs/aphrodite.cron
MAILTO='william@tuffbizz.com'

0 12 * * * /home/william/scripts/cron-shunt /home/william/scripts/dailymail.rb
