# scripts/crontabs/athena.cron
0 0,12 * * * /home/william/scripts/cron-shunt twitter-fortune-bot.pl
