# scripts/crontabs/athena.cron
0,30 * * * * /home/william/scripts/cron-shunt twitter-fortune-bot.pl
