# scripts/crontabs/athena.cron
0 0,12 * * * /home/william/bots/twitter-fortune-bot.pl
