# scripts/crontabs/athena.cron
MAILTO='william@tuffbizz.com'

0,30 * * * * /home/william/scripts/cron-shunt twitter-fortune-bot.pl
