# scripts/crontabs/mars.cron
PATH=/usr/local/bin:/usr/bin:/bin:/Users/william/scripts

0 12 * * * /Users/william/scripts/cron-shunt dailymail.rb
