# scripts/crontabs/mars.cron
MAILTO='william@tuffbizz.com'

0 12 * * * /Users/william/scripts/cron-shunt /Users/william/scripts/dailymail.rb
